# config/files/etc/systemd/system/tailscaled.service.d/override.conf
[Service]
ExecStartPre=nft -f '/usr/mullvad-ts.nft'
ExecStopPost=nft -f '/usr/mullvad-ts-cleanup.nft'